# supervisord.conf
[program:dropbear]
command=/usr/sbin/dropbear -RFEwg -p 22